# Nginx-web-server/Nginx-LB01-server/conf/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
>/etc/udev/rules.d/70-persistent-net.rules

## npd service start @root @2018-3-14 @cr001
#/etc/init.d/ntpd start #( At heartbeat up)

## nginx and php start  @root @2018-3-14 @cr002
#/application/nginx/sbin/nginx #( At heartbeat up)
